# crontab/crontab_file
# crontab schedule file

# UTC: 18:08 for Eastern Servers, Beijing: 02:08
# UTC: 05:08 for Western Servers
# minute hour day month week cmd
  8      18   *   *     *    /dinstar/cloudenv/crontab/crontab_task_day.sh

# minute hour day month week cmd
  */1    *    *   *     *    /dinstar/cloudenv/crontab/crontab_task_minute.sh

